# server_proxy/scripts/proxycronjob.cron
*/5 * * * * /var/imaproxy/scripts/proxycronjob

