# server_proxy/scripts/proxycronjob.cron
*/5 * * * * ${IMA_PROXY_DATA_PATH}/scripts/proxycronjob

