# server_proxy_br/scripts/bridgecronjob.cron
*/5 * * * * ${IMA_BRIDGE_DATA_PATH}/scripts/bridgecronjob

